# music/robot_music.service
[Unit]
Description=MVRT Robot Music

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/robot_music.sh

[Install]
WantedBy=multi-user.target